# src/servers/dev1/dev1.crontab
# dev1 Crontab
# NOTE: Do not store environment variables in this file!
# They will be sent via email if a cron job fails which could expose sensitive creds
# Instead, place them in `~/.env-crontab` and source that file on
# cron jobs that require environment variables


# GLOBAL SETUP
PATH=/home/linuxbrew/.linuxbrew/bin:/home/linuxbrew/.linuxbrew/sbin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
MAILTO=admin@shieldmanagement.net


# Alchemist (DO NOT use `--greedy` here since it will break Docker in prod)
10 2 * * * cronic alchemist --update
30 2 * * * cronic alchemist --backup --force
45 2 * * * cronic cp "$HOME/alchemist/backup/Brewfile" "$HOME/.dotfiles/src/$HOSTNAME/Brewfile"

# Docker
0 3 * * 0 cronic docker system prune -f

# Healthchecks
*/3 * * * * cd "$HOME"/git/server-infra && cronic ./healthchecks/docker-healthcheck.sh

# GitHub Archive
0 1 * * * . "$HOME"/.env-crontab; . "$HOME/.keychain/$HOSTNAME-sh"; cronic venv/github-archive/bin/github-archive --users justintime50 --orgs "$GITHUB_ARCHIVE_ORGS" --gists justintime50 --clone --pull --forks --token "$GITHUB_TOKEN" --log_level warning

# Pullbug
0 9 * * 1 . "$HOME"/.env-crontab; cronic venv/pullbug/bin/pullbug --pulls --issues --drafts --github_token "$GITHUB_TOKEN" --github_context users --github_owner Justintime50 --discord --discord_url "$PULLBUG_DISCORD_URL" --log_level warning

# Unifi Config Backups
0 4 * * * . "$HOME"/.env-crontab; UNIFI_USERNAME="$UNIFI_USERNAME" UNIFI_PASSWORD="$UNIFI_PASSWORD" UNIFI_BACKUP_LOCATION="$UNIFI_BACKUP_LOCATION" cronic venv/get_last_unifi_console_backup/bin/python "$HOME"/git/os-scripting/src/ubios/backups/get_last_console_backup.py
5 4 * * 0 . "$HOME"/.env-crontab; cronic find "$UNIFI_BACKUP_LOCATION" -mtime +90 -delete
