# src/servers/dev1/dev1.crontab
# dev1 Crontab
# NOTE: Do not store environment variables in this file!
# They will be sent via email if a cron job fails which could expose sensitive creds
# Instead, place them in `~/.env-crontab` and source that file on
# cron jobs that require environment variables


# GLOBAL SETUP
PATH=/home/linuxbrew/.linuxbrew/bin:/home/linuxbrew/.linuxbrew/sbin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
MAILTO=admin@shieldmanagement.net


# Alchemist (DO NOT use `--greedy` here since it will break Docker in prod)
10 9 * * * cronic alchemist --update
30 9 * * * cronic alchemist --backup --force
45 9 * * * cronic cp "$HOME/alchemist/backup/Brewfile" "$HOME/.dotfiles/src/servers/$(hostname)/Brewfile"

# Docker
0 10 * * 0 cronic docker system prune -f
5 10 * * * cronic docker exec -t $(docker ps --filter "name=hub_hub.1" --format "{{.ID}}" | head -n 1) /opt/registry-ui -purge-tags

# Healthchecks
*/3 * * * * cd "$HOME"/git/server-infra && cronic ./scripts/healthchecks/docker-healthcheck.sh
*/5 * * * * . "$HOME"/.env-crontab; export GOVEE_API_KEY; cronic venv/govee_smartplug_healthcheck/bin/python "$HOME"/git/server-infra/scripts/healthchecks/govee_smartplug_healthcheck.py

# GitHub Archive
0 8 * * * . "$HOME"/.env-crontab; . "$HOME/.keychain/$(hostname)-sh"; cronic venv/github-archive/bin/github-archive --users justintime50 --orgs "$GITHUB_ARCHIVE_ORGS" --gists justintime50 --clone --pull --forks --token "$GITHUB_TOKEN" --log_level warning

# Pullbug
0 16 * * 1 . "$HOME"/.env-crontab; cronic venv/pullbug/bin/pullbug --pulls --issues --drafts --github_token "$GITHUB_TOKEN" --github_context users --github_owner Justintime50 --discord --discord_url "$PULLBUG_DISCORD_URL" --log_level warning

# Unifi Config Backups
0 11 * * * . "$HOME/.keychain/$(hostname)-sh"; . "$HOME"/.env-crontab; export UNIFI_USERNAME UNIFI_PASSWORD UNIFI_BACKUP_LOCATION; cd "$HOME"/git/server-infra/scripts/backups && cronic ./unifi-backup-maintenance.sh
